# tetra_monitor.service
[Unit]
Description=Pi Eye TETRA Monitor Dashboard
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/tetra/tetra_dashboard.py
WorkingDirectory=/opt/tetra
Restart=on-failure
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1

[Install]
WantedBy=multi-user.target
